# nacc_manager/center-web-manager 8008.ini
[program:center-web-manager-8008] 
directory=/yihnac/center-web-manager
command=python3 manage.py runserver 0.0.0.0:8008
redirect_stderr=true
stdout_logfile=/yihnac/supervisor/logs/8008_stdout.out
priority=991  ;进程启动优先级，默认999，值小的优先启动
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root